# ms_jk_task/instance_tpl/supervisord.conf
[supervisord]
logfile=/workspace/instance/log/supervisord.log ; main log file; default $CWD/supervisord.log
nodaemon = true
[program:gunicorn]
directory=/workspace
command=gunicorn wsgi:flask_app -c instance/gunicorn.conf.py
stdout_logfile=/workspace/instance/log/gunicorn.stdout
stderr_logfile=/workspace/instance/log/gunicorn.stderr
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
